# systemd/myclock-stop.service
[Unit]
Description="Service to stop myclock.service"

[Service]
ExecStart=systemctl stop myclock


