# systemd/myclock.service
[Unit]
Description=Start myclock service
#Requires=gw_log.service
#After=gw_log.service
DefaultDependencies=no

[Service]
Type=simple
Restart=on-failure
RestartSec=5

# set User to find locally installed python modules
User=pi

# Commands to run before starting service
ExecStartPre=-cp -b /home/pi/myprojects/myclock/myclock.stdout.log \
		/home/pi/myprojects/myclock/myclock.stdout.old

ExecStartPre=-cp -b /home/pi/myprojects/myclock/myclock.stderr.log \
		/home/pi/myprojects/myclock/myclock.stderr.old

ExecStartPre=-rm /home/pi/myprojects/myclock/myclock.stdout.log
ExecStartPre=-rm /home/pi/myprojects/myclock/myclock.stderr.log

# Start service
# ExecStart for non virtual env
#ExecStart=/usr/bin/python -u\
#	 /home/pi/myprojects/myclock/myclock.py -b 0 

# ExecStart for virtual env
ExecStart=/home/pi/.virtualenvs/myclock/bin/python -u\
	 /home/pi/myprojects/myclock/myclock.py -b 0 

StandardOutput=file:/home/pi/myprojects/myclock/myclock.stdout.log
StandardError=file:/home/pi/myprojects/myclock/myclock.stderr.log

[Install]
WantedBy=multi-user.target
